# doc/system/crontab
* * * * * cd /home/pi/rtjam && /home/pi/rtjam/update-rtjam.bash
